# flaskblog.conf
[program:flaskblog]
command = gunicorn -b wsgi_gunicorn:api
directory = /chatroom
user = root